# configs/overlay/etc/inittab
::sysinit:/etc/init.d/rc
::askfirst:/etc/initshell
::restart:/sbin/init
::ctrlaltdel:/bin/busybox reboot -f
::shutdown:umount -a -r